# src/minecraftd-save@.timer
[Unit]
Description=Daily save on Minecraft server %I

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target