# src/minecraftd-restart@.service
[Unit]
Description=Restart %I Minecraft server

[Service]
Type=oneshot
ExecStart=/usr/bin/minecraftctl restart %i