# src/minecraftd-backup@.service
[Unit]
Description=Backup %I Minecraft server

[Service]
Type=oneshot
ExecStart=/usr/bin/minecraftctl backup %i