# src/minecraftd-save@.service
[Unit]
Description=Save %I Minecraft server

[Service]
Type=oneshot
ExecStart=/usr/bin/minecraftctl save %i